# hub/systemd/rclone-catchup.service
[Unit]
Description=Catch up from R2 on boot
After=network-online.target
Before=rclone-webdav.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/rclone sync cloud-raw:e2ee-sync /data/encrypted \
  --checksum \
  --transfers 8 \
  --verbose \
  --log-file /var/log/rclone-catchup.log

[Install]
WantedBy=multi-user.target
